# docker/supervisord.conf
[program:php-fpm]
command=/usr/sbin/php-fpm -F -R -c /etc/php/8.2/fpm/php-fpm.conf
autostart=true
autorestart=true
startretries=3
stopasgroup=true